# .config/systemd/user/hyper-push.service
[Unit]
Description=Sync dotfiles to GitHub

[Service]
Type=oneshot
ExecStart=/home/duane/local/bin/hyper-push
