# .config/systemd/user/vdirsyncer.service
[Unit]
Description=Synchronize calendars with vdirsyncer

[Service]
Type=oneshot
ExecStart=/usr/bin/vdirsyncer sync
