# .config/systemd/user/mpd-mpris.service.d/override.conf
[Service]
ExecStart=
ExecStart=mpd-mpris -no-instance -host 127.0.0.1 -port 6601